# overlays/System/Library/LaunchDaemons/com.apple.getty.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.apple.getty.plist

  Apple-faithful single-plist getty service. Shape matches Apple's
  /System/Library/LaunchDaemons/com.apple.getty.plist on macOS:

      ProgramArguments = ("/usr/libexec/getty", "std.9600", "console")
      KeepAlive = true
      RunAtLoad = true (implicit in com.apple.getty; we set it
                       explicitly so first-boot fires getty without
                       waiting on an external trigger)

  Why /dev/console (not ttyu0 / ttyv0):
    On both Darwin and FreeBSD, /dev/console is a kernel-provided
    character device routed to whichever console driver is active for
    this boot — UART on serial-only boots, vt(4) on VGA, both when
    boot.loader.conf says "Dual Console". The launchd job doesn't have
    to know which physical device is behind it; the kernel handles
    the multiplexing. This eliminates the missing-device respawn-loop
    problem and removes the need for the per-device wrapper script
    that the earlier two-plist approach used.

  Why std.9600 (gettytab class):
    Matches the exact string Apple uses. FreeBSD's /etc/gettytab
    already ships this class (BSD-inherited). The class governs baud
    rate, terminal type defaults (vt100 for std.* classes), and
    parity. A user wanting a richer TERM (xterm-color, cons25) can
    override after login.

  We keep using FreeBSD's /usr/libexec/getty for now (BSD-derived,
  works on /dev/console identically to Apple's getty). Porting
  Apple's getty itself is deferred until after the broader pkgbase
  cleanup pass per the user roadmap.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apple.getty</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/libexec/getty</string>
        <string>std.9600</string>
        <string>console</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
