# overlays/System/Library/LaunchDaemons/org.freebsd.getty.vty0.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.getty.vty0.plist

  Spawns getty on /dev/ttyv0 — FreeBSD's first VGA virtual terminal,
  served by the vt(4) console driver (vt_efifb on UEFI hardware,
  vt_vga on legacy BIOS). This is the console real-hardware boots
  see; serial-console-only boots (qemu with -serial stdio, no VGA)
  don't have /dev/ttyv0 at all.

  Wrapped via launchd-getty-wrapper so missing-device hardware
  (serial-only) doesn't trigger a launchd respawn loop — the wrapper
  checks /dev/ttyv0 existence first and sleeps for 24h if absent.

  KeepAlive=true: respawn after each login session ends.

  Pc = the gettytab entry for PC keyboards (covers most VGA console
  setups; real hardware, VirtualBox VGA, VMware VGA, etc.).

  Ported from sibling sockets-version freebsd-launchd repo.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.getty.vty0</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/libexec/launchd-getty-wrapper</string>
        <string>Pc</string>
        <string>ttyv0</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
