# overlays/System/Library/LaunchDaemons/org.freebsd.getty.console.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.freebsd.getty.console.plist

  Spawns getty on /dev/ttyu0 (the serial console — boot.loader.conf
  routes console=comconsole). Equivalent to the /etc/ttys line a
  conventional FreeBSD boot would have:

      ttyu0  "/usr/libexec/getty 3wire"  vt100  on  secure

  KeepAlive=true so launchd respawns getty after each login session,
  matching init's classic behavior. RunAtLoad=true so getty fires
  immediately on launchd startup — without it, the boot test would
  never see "login:".

  Wrapped via launchd-getty-wrapper so missing-device hardware
  (video-only) doesn't trigger a launchd respawn loop — the wrapper
  checks /dev/ttyu0 existence first and sleeps for 24h if absent.

  No StandardInPath/Out/ErrPath here on purpose: getty opens its
  controlling tty itself with TIOCSCTTY. If launchd pre-opens those
  fds and execs into getty with them inherited, the controlling-tty
  acquisition silently fails and login never appears.

  Ported from sibling sockets-version freebsd-launchd repo, which
  proved this exact shape boots to login.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.freebsd.getty.console</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/libexec/launchd-getty-wrapper</string>
        <string>3wire</string>
        <string>ttyu0</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
